# service/altcoin-um/data-source-altcoin-um.timer
[Unit]
Description=Run data-source-altcoin-um daily at 08:20 UTC

[Timer]
OnCalendar=*-*-* 08:20:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
