# service/altcoin-spot/data-source-altcoin-spot.timer
[Unit]
Description=Run data-source-altcoin-spot daily at 08:00 UTC

[Timer]
OnCalendar=*-*-* 08:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
